# xfce/.config/autostart/autostart.desktop
[Desktop Entry]
Type=Application
Name=autostart
Exec=autostart
StartupNotify=false
Terminal=false
