# autostart/.config/autostart/fcitx5.desktop
[Desktop Entry]
Name=fcitx5
Comment=fcitx5
Exec=fcitx5
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
